# server/deploy/systemd/portal.service
[Unit]
Description=Reverse SSH Tunnel Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PORTAL_DIR__
EnvironmentFile=__PORTAL_DIR__/.env
ExecStart=__PORTAL_DIR__/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port __PORT__ --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
